# debian/qemu-kvm-extras-static.postrm
#!/bin/sh

set -e

# apply /etc/sysctl.d settings
if [ "$1" = purge ]; then
    invoke-rc.d procps start
fi

#DEBHELPER#
